# GRFICSv3/simulation/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:simulation]
command=/app/simulation/simulation
autorestart=true
stdout_logfile=/var/log/supervisor/simulation.log
stderr_logfile=/var/log/supervisor/simulation.err

[program:modbus]
directory=/app/modbus
command=/app/modbus/run_all.sh
autorestart=true
; Python buffers stdout/stderr when not a TTY — supervisor redirects to files, so without
; this, modbus.log/modbus.err stay empty until buffer fills or the process exits.
environment=PYTHONUNBUFFERED=1
stdout_logfile=/var/log/supervisor/modbus.log
stderr_logfile=/var/log/supervisor/modbus.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
